# docs/operations/pepepow-explorer.service
[Unit]
Description=PEPEPOW eIquidus Explorer
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/explorer
Environment=NODE_ENV=production
Environment=PATH=/home/ubuntu/.nvm/versions/node/v22.11.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/ubuntu/.nvm/versions/node/v22.11.0/bin/npm start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
